# portal/PortalSvr/supervisord.conf
[supervisord]
nodaemon=true

[program:PortalSvr]
command=dotnet PortalSvr.dll

[program:sendmail]
command=/etc/init.d/sendmail start
